# lib/systemd/system/powermoderd.service
[Unit]
Description=powermdoer daemon

[Service]
Type=simple
ExecStart=bash -ic "powermoderd init"

[Install]
WantedBy=multi-user.target
